# node_server.service
[Unit]
Description=Node Server Service

[Service]
ExecStart=/usr/bin/node /home/pi/server.js
Restart=always

WantedBy=multi-user.target